# server/deploy/supervisord_conf.ini
[program:ext_newsletter]
user=root
command=gunicorn -c deploy/prd_newsletter.py run:app
autostart=true
autorestart=true
stopasgroup=true
directory=/var/www/ext_newsletter/server